# minio.service
[Unit]
Description=MinIO Object Storage
After=network.target

[Service]
Type=simple
User=root
Environment="MINIO_ROOT_USER=heytex_admin"
Environment="MINIO_ROOT_PASSWORD=heytex_minio_2024"
ExecStart=/usr/local/bin/minio server /root/heytex/data/minio --console-address ":9001"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
